# deploy/systemd/rmuc-live-sync.service
[Unit]
Description=RMUC live schedule and mini-program prediction sync
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# Change User/Group to the deployment user that owns /opt/douququ.
User=douququ
Group=douququ
WorkingDirectory=/opt/douququ
Environment=PYTHONPATH=/opt/douququ
Environment=RMUC_MINI_PROGRAM_ENABLED=1
ExecStart=/opt/douququ/.venv/bin/python /opt/douququ/scripts/sync_rmuc_live.py --mini-program-ttl-seconds 300 --mini-program-refresh-window-seconds 60
